# .config/systemd/user/xsrv-backup.service
[Unit]
Description=Xserver Backup Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/xsrv-backup.sh

