# cron/crontab
# 2026-04-06 update
# Edit this file to introduce tasks to be run by cron.
#
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
#
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
#
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
#
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
#
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
#
# For more information see the manual pages of crontab(5) and cron(8)
#
# m h  dom mon dow   command
# Gospel-Haiku Passwd marge & backup（毎日23:40）
40 23 * * * /usr/local/bin/automerge.sh >> /tmp/cron.log 2>&1
# GH+minorugh.com dotfils git auto push （毎日23:50）
50 23 * * * /usr/local/bin/autobackup.sh >> /tmp/cron.log 2>&1
# Emacs ゴミ箱をシステムゴミ箱へ移送（毎月1日 23:45）
45 23 1 * * /usr/local/bin/emacs-trash-sweep.sh >> /tmp/cron.log 2>&1
# xmodmap 毎分再適用（失速対策）
* * * * * DISPLAY=:0 xmodmap ~/.Xmodmap 2>> /tmp/cron.log

