# cron/crontab
40 23 * * * /usr/local/bin/automerge.sh >> /tmp/cron.log 2>&1
50 23 * * * /usr/local/bin/autobackup.sh >> /tmp/cron.log 2>&1
* * * * * DISPLAY=:0 /usr/bin/xmodmap /home/minoru/.Xmodmap >> /tmp/xmodmap.log 2>&1
# Mattermost DBバックアップ（毎日午前3時）
0 3 * * * /usr/local/bin/mattermost-backup.sh >> /tmp/cron.log 2>&1
